# systemd/azazel-epd.service.d/10-portal-detect.conf
[Unit]
Wants=network-online.target
After=network-online.target dhcpcd.service

[Service]
EnvironmentFile=/etc/default/azazel-zero
ExecStartPost=/usr/local/bin/portal_detect.sh "$OUTIF"