# systemd/azazel-boot-splash.service
# /etc/systemd/system/azazel-boot-splash.service（最小形）
[Unit]
Description=Azazel-Zero E-Paper Boot Splash
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=azazel
Environment=PYTHONPATH=/opt/waveshare-epd/RaspberryPi_JetsonNano/python:/opt/waveshare-epd/RaspberryPi_JetsonNano/python/lib
ExecStart=/usr/bin/python3 /home/azazel/Azazel-Zero/py/boot_splash_epd.py
WorkingDirectory=/home/azazel/Azazel-Zero/py
Restart=on-failure

[Install]
WantedBy=multi-user.target