# systemd/azazel-epd.service
[Unit]
Description=Azazel EPD Boot/Shutdown Animations
DefaultDependencies=no
Before=graphical.target shutdown.target
After=local-fs.target
ConditionPathExists=/home/azazel/Azazel-Zero/py/boot_splash_epd.py

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/default/azazel-zero
Environment=EPD_PY=/home/azazel/Azazel-Zero/py/boot_splash_epd.py
ExecStart=/usr/bin/python3 $EPD_PY --mode start --steps 10 --frame-sec 0.25
ExecStartPost=/usr/bin/python3 $EPD_PY --mode info
ExecStop=/usr/bin/python3 $EPD_PY --mode shutdown
TimeoutStopSec=12

[Install]
WantedBy=multi-user.target