# systemd/azazel-console.service
[Unit]
Description=Azazel-Zero Console (tmux session)
After=network-online.target azazel-epd.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=azazel
EnvironmentFile=-/etc/default/azazel-zero
Environment=HOME=/home/azazel
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin
ExecStart=/bin/bash -lc '/usr/local/bin/azazel_console.sh'
ExecStop=/bin/bash -lc '/usr/bin/tmux has-session -t azazel 2>/dev/null && /usr/bin/tmux kill-session -t azazel || true'

[Install]
WantedBy=multi-user.target